# backend/server/php/profit_calculation.cron
# daily currency/commodities exchange proposals calculation and logging
0 9 * * * /usr/local/bin/php /var/www/html/scripts/profitCalculation.php >> /var/log/profit_calculation.log
# An empty line is required at the end of this file for a valid cron file.